# deploy/systemd/coletor.timer
[Unit]
Description=agente-tjms: coletor diário (04:30)

[Timer]
OnCalendar=*-*-* 04:30:00
RandomizedDelaySec=10min
Persistent=true

[Install]
WantedBy=timers.target
