# deploy/systemd/rastreador.timer
[Unit]
Description=agente-tjms: rastreador diário (21:00)

[Timer]
OnCalendar=*-*-* 21:00:00
RandomizedDelaySec=10min
Persistent=true

[Install]
WantedBy=timers.target
