# deploy/systemd/rastreador.service
[Unit]
Description=agente-tjms: rastreador diário de acórdãos

[Service]
Type=oneshot
WorkingDirectory=%h/projetos/agente-tjms
ExecStart=%h/projetos/agente-tjms/.venv/bin/agente-tjms rastrear-acordaos
